# packages/wolkenbruch/wolkenbruch-0.4.3.tar.gz/wolkenbruch-0.4.3/extra/systemd/wolkenbruch@.service
[Unit]
Description = Check whether it’s likely to rain and send an e-mail reminder for user %i

[Service]
ExecStart = /usr/bin/wolkenbruch
User = %i
